# deploy/glanceos-backend.service
[Unit]
Description=GlanceOS Backend
After=network.target

[Service]
ExecStart=/home/pi/glanceos/backend/venv/bin/python /home/pi/glanceos/backend/run.py
WorkingDirectory=/home/pi/glanceos/backend
Restart=always
RestartSec=5
User=pi
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
